# dbus/blivet.service
[Unit]
Description=Blivet Service

[Service]
Type=dbus
BusName=org.blivet.Blivet
ExecStart=/usr/libexec/blivetd
